# crontab
2 * * * * root doAction >> /var/log/cron.log
