# crontab
2 * * * * root node action.js >> /var/log/cron.log
